# usr/local/etc/supervisor.d/dnsmasq.ini
[program:dnsmasq]
command = /usr/local/sbin/dnsmasq --keep-in-foreground -C /usr/local/etc/dnsmasq.conf

autostart=true
autorestart=true
startsecs=3
startretries=3

stderr_logfile=/usr/local/var/log/dnsmasq.log
stdout_logfile=/usr/local/var/log/dnsmasq.log

user=root
