# usr/local/etc/supervisor.d/nginx.ini
[program:nginx]
command =  /usr/local/bin/nginx -g "daemon off;"

autostart=true
autorestart=true
startsecs=3
startretries=3

stdout_logfile=/usr/local/var/log/nginx/access.log
stderr_logfile=/usr/local/var/log/nginx/error.log

user=yaccai
