# usr/local/etc/supervisor.d/watch.ini
[program:watch]
command = /usr/local/bin/watchman --foreground --logfile=/usr/local/var/run/watchman/log --log-level=1 --sockname=/usr/local/var/run/watchman/sock --statefile=/usr/local/var/run/watchman/state --pidfile=/usr/local/var/run/watchman/pid

autostart=true
autorestart=true
startsecs=3
startretries=3

stderr_logfile=/usr/local/var/run/watchman/log
stdout_logfile=/usr/local/var/run/watchman/log

user=yaccai

